# radius/upstart/foxpass-radius-agent.conf
description "Foxpass RADIUS agent"
author  "Aren Sandersen <aren@foxpass.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir /opt/bin
exec python /opt/bin/agent.py
respawn
