# scripts/rc.local
echo 1 > /proc/sys/net/ipv4/ip_forward
exit 0
